# scripts/init/joblens.service.in
[Unit]
Description=JobLens - high-performance job monitoring daemon
Documentation=https://github.com/joblens-org/JobLens
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/JobLens
@JOBLENS_SERVICE_ENV@
ExecStartPre=@CMAKE_INSTALL_FULL_BINDIR@/JobLens -v
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/JobLens \
          -m service \
          -c /etc/JobLens/config.yaml
ExecStopPost=-/bin/rm -f /var/JobLens/JobLens.lock
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=joblens

[Install]
WantedBy=multi-user.target
